# marathon_docker/dummy_chimpcheck.conf

[program:startdocker]
command=echo "Hello"