# docker/chimpcheck.conf

[program:startdocker]
command=bash /chimpcheck/ChimpWebDemo/start_docker.bash